# templates/gerrit.service.erb
# WARNING THIS FILE IS MANAGED BY PUPPET
[Unit]
Description=Gerrit Code Review
After=network.target

[Service]
EnvironmentFile=/etc/default/gerritcodereview
ExecStart=/usr/bin/java $JAVA_OPTIONS -jar <%= @gerrit_home %>/bin/gerrit.war daemon -d $GERRIT_SITE --run-id `date +%s`.$$
User=<%= @gerrit_user %>
Group=<%= @gerrit_user %>

[Install]
WantedBy=multi-user.target
